# src/init.d/torrus.in
#!@SHELL@
#
# init.d script for Torrus
# Install it as /etc/init.d/torrus (most UNIXes),
# or as /usr/local/etc/rc.d/torrus.sh (FreeBSD), or probably somewhere else.
#
# Stanislav Sinyagin <ssinyagin@k-open.com>
#
#
### chkconfig info
# chkconfig: 2345 90 10
# description: Starts/Stops Torrus collectors and monitors
###
#
### BEGIN INIT INFO
# Provides:          torrus-collector torrus-monitor
# Required-Start:    $local_fs $remote_fs $syslog $named $network $time
# Required-Stop:     $local_fs $remote_fs $syslog $named $network
# Should-Start:
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Torrus collectors and monitors
# Description:       Start/stop Torrus collectors and monitors
### END INIT INFO

prefix=@prefix@
sysconfdir=@sysconfdir@
pkghome=@pkghome@
cmddir=@pkgbindir@
piddir=@piddir@
sitedir=@sitedir@
torrususer=@torrus_user@

. @cfgdefdir@/initscript.conf
if test -f @siteconfdir@/initscript.siteconf; then
  . @siteconfdir@/initscript.siteconf
fi

if test ! -d ${piddir}; then
  install --owner ${torrususer} --directory ${piddir}
fi

if test `ls -ld ${piddir} | awk 'NR==1 {print $3}'` != ${torrususer}; then
  chown ${torrususer} ${piddir}
fi

START_ARGS="--cmd=start"
STOP_ARGS="--cmd=stop"

if test x"${TORRUS_USER}" != x; then
  START_ARGS=${START_ARGS}' --user='${TORRUS_USER}
  STOP_ARGS=${STOP_ARGS}' --user='${TORRUS_USER}
fi

if test x"${TORRUS_COLLECTOR_CMDOPTS}" != x; then
  START_ARGS=${START_ARGS}' --copts="'${TORRUS_COLLECTOR_CMDOPTS}'"'
fi

if test x"${TORRUS_MONITOR_CMDOPTS}" != x; then
  START_ARGS=${START_ARGS}' --mopts="'${TORRUS_MONITOR_CMDOPTS}'"'
fi

if test x"${TORRUS_STOP_TIMEOUT}" != x; then
  STOP_ARGS=${STOP_ARGS}' --timeout='${TORRUS_STOP_TIMEOUT}
fi


start_daemons () {
  ${cmddir}/launcher ${START_ARGS} 2>&1 | tee -a @logdir@/launcher_start.log

  # RHEL based systems (RHEL, CentOS, Fedora) ignore the KXXtorrus script
  # unless the corresponding lock is present
  if test -d /var/lock/subsys; then
    touch /var/lock/subsys/torrus
  fi
}

stop_daemons () {
  ${cmddir}/launcher ${STOP_ARGS} 2>&1 | tee -a @logdir@/launcher_stop.log

  # RHEL specifics
  if test -d /var/lock/subsys; then
    rm -f /var/lock/subsys/torrus
  fi
}

case "$1" in
  'start')    start_daemons
  ;;

  'stop')     stop_daemons
  ;;
  
  'restart')  stop_daemons; start_daemons
  ;;

  'status')   ${cmddir}/launcher --cmd=status
  ;;
  
  *)  echo "Usage: $0 [start|stop|restart]"
  ;;
  esac

# Local Variables:
# mode: shell-script
# sh-shell: sh
# indent-tabs-mode: nil
# sh-basic-offset: 2
# End:
